# conf/supervisord_uwsgi.ini
[program:redbutterfly]
command=/usr/sbin/uwsgi -x /home/hanfei/project/django/redbutterfly/redbutterfly_blog/bin/django_redbutterfly.xml --buffer-size 32768
directory=/home/hanfei/project/django/redbutterfly/redbutterfly_blog/bin
user=root
startsecs=0
stopwaitsecs=0
autostart=true
autorestart=true
stdout_logfile=/home/hanfei/project/django/redbutterfly/redbutterfly_blog/log/redbutterfly.stdout.log
stderr_logfile=/home/hanfei/project/django/redbutterfly/redbutterfly_blog/log/redbutterfly.stderr.log
stopsignal=QUIT
killasgroup=true
