# conf/supervisord_gunicorn.ini
[program:redbutterfly]
command=gunicorn conf.wsgi:application -w 4 -b 127.0.0.1:9001 -k gevent --max-requests 500
directory=/home/hanfei/project/django/redbutterfly/redbutterfly_blog
user=root
startsecs=0
stopwaitsecs=0
autostart=true
autorestart=true
stdout_logfile=/home/hanfei/project/django/redbutterfly/redbutterfly_blog/log/redbutterfly.stdout.log
stderr_logfile=/home/hanfei/project/django/redbutterfly/redbutterfly_blog/log/redbutterfly.stderr.log
stopsignal=QUIT
killasgroup=true
